# ansible/roles/post_upgrade_reboot/templates/post-upgrade-reboot.service.j2
# {{ ansible_managed }}
[Unit]
Description=Reboot after package upgrades if required (quorum-aware)
Documentation=file:/usr/local/sbin/post-upgrade-reboot.sh

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/post-upgrade-reboot.sh
# Worst case: drain ({{ post_upgrade_reboot_drain_timeout_seconds }}s) + etcd
# snapshot (90s) + 90s wait for a stuck PID 1 before reboot -f, plus margin.
TimeoutStartSec=25min
